# com.brewery.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Build:   cd agent && go build -o brewery-agent .
  Install: sudo cp agent/brewery-agent /usr/local/bin/brewery-agent
           sudo cp com.brewery.agent.plist /Library/LaunchDaemons/com.brewery.agent.plist
           sudo launchctl load /Library/LaunchDaemons/com.brewery.agent.plist
  Unload:  sudo launchctl unload /Library/LaunchDaemons/com.brewery.agent.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brewery.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/brewery-agent</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BREWERY_SERVER_URL</key>
    <string>http://your-brewery-server:6502</string>
    <!-- Leave BREWERY_API_KEY empty if the server has no API key set -->
    <key>BREWERY_API_KEY</key>
    <string></string>
    <!-- brew must be on PATH; include both ARM and Intel prefixes -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Run every 15 minutes -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/brewery-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/brewery-agent.log</string>
</dict>
</plist>
